# cpu-bot.service
[Unit]
Description=CPU Keepalive Discord Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/discord-bot24
ExecStart=/usr/bin/python3 cpu_bot.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

StandardOutput=journal
StandardError=journal
SyslogIdentifier=cpu-bot

[Install]
WantedBy=multi-user.target
