# oracle-monitor.service
[Unit]
Description=HomeServer Monitor Discord Bot
After=network.target

[Service]
Type=simple
User=homeserver
WorkingDirectory=/home/homeserver/discord-bot24
ExecStart=/usr/bin/python3 bot.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

StandardOutput=journal
StandardError=journal
SyslogIdentifier=homeserver-monitor-bot

[Install]
WantedBy=multi-user.target
